# factory-config/config/includes.chroot/etc/cron.daily/anodyne-housekeeper
#!/bin/sh
/usr/bin/anodyne-housekeeper.sh
